# qubes-vm/qrexec/files/qubes-qrexec-agent.initd
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="qrexec agent"
description="Qubes remote exec agent"
command=/usr/lib/qubes/qrexec-agent
command_background=1
pidfile=/var/run/qrexec-agent.pid
#command_args="${qrexec_args}"
output_log=/var/log/qubes-qrexec.log

depend() {
	after qubes-db
}
