# qubes-vm/core-agent/files/qubes-services.initd
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="qubes-sysinit"
description="Init Qubes Services settings"
depend() {
	need qubes-db
}

start() {
	/usr/lib/qubes/init/qubes-sysinit.sh 1>/dev/null 2>&1
	exit $?
}
